# files/etc/systemd/system/docker.service.d/docker-storage.conf
[Service]
EnvironmentFile=-/etc/sysconfig/docker-storage
ExecStart=
ExecStart=/usr/bin/dockerd $DOCKER_STORAGE_OPTIONS
