# files/etc/systemd/system/docker.service.d/mountflags-slave.conf
[Service]
# MountFlags "slave" helps to prevent "device busy" errors on CentOS/RedHat 7.3 kernels
MountFlags=slave
